# supervisord-sshd.conf
[program:sshd]
command=/start-sshd.wrapper.sh
numprocs=1
autostart=true
autorestart=true
